# packaging/control-ofc-daemon.install
post_install() {
    echo ":: Enable the daemon:"
    echo "   sudo systemctl enable --now control-ofc-daemon"
    echo ""
    echo ":: Load sensor modules now (loaded automatically on future boots):"
    echo "   sudo systemctl start systemd-modules-load"
    echo ""
    echo ":: Next steps:"
    echo "   - Install the GUI:  yay -S control-ofc-gui"
    echo "   - AMD GPU fans require kernel parameter: amdgpu.ppfeaturemask=0xffffffff"
    echo "   - Full setup guide: /usr/share/doc/control-ofc-daemon/"
}

post_upgrade() {
    echo ":: Restart to apply changes:"
    echo "   sudo systemctl restart control-ofc-daemon"

    _strip_legacy_runtime_sections
}

# Auto-strip [profiles] and [startup] sections from /etc/control-ofc/daemon.toml
# if they exist. Back up the original to daemon.toml.pre-1.1.2.bak first. Safe
# to re-run — no-op when the sections are already gone. Uses awk so we do not
# depend on sed GNU extensions.
_strip_legacy_runtime_sections() {
    local config="/etc/control-ofc/daemon.toml"
    [ -f "$config" ] || return 0

    # Grep for top-of-line section headers only; ignore comments and dotted keys.
    if ! grep -qE '^[[:space:]]*\[(profiles|startup)\][[:space:]]*$' "$config"; then
        return 0
    fi

    local backup="${config}.pre-1.1.2.bak"
    if ! cp -a "$config" "$backup" 2>/dev/null; then
        echo ":: WARNING: could not back up $config — leaving legacy sections in place."
        echo "   Remove [profiles] and [startup] manually before upgrading to 1.2.0."
        return 0
    fi

    local tmp="${config}.tmp-$$"
    awk '
        BEGIN { skip = 0 }
        /^[[:space:]]*\[profiles\][[:space:]]*$/ { skip = 1; next }
        /^[[:space:]]*\[startup\][[:space:]]*$/  { skip = 1; next }
        /^[[:space:]]*\[/                         { skip = 0 }
        skip == 0 { print }
    ' "$config" > "$tmp" || {
        echo ":: WARNING: failed to strip legacy sections from $config."
        echo "   Backup saved to $backup. Remove [profiles] and [startup] manually."
        rm -f "$tmp"
        return 0
    }

    # Preserve mode/owner from the original (should be root:root 0644 from PKGBUILD).
    chmod --reference="$config" "$tmp" 2>/dev/null || chmod 0644 "$tmp"
    chown --reference="$config" "$tmp" 2>/dev/null || true

    if mv -f "$tmp" "$config"; then
        echo ":: Stripped legacy [profiles] / [startup] sections from $config."
        echo "   Original saved as $backup."
    else
        echo ":: WARNING: failed to replace $config — backup is at $backup."
        rm -f "$tmp"
    fi
}

pre_remove() {
    if systemctl is-active --quiet control-ofc-daemon.service; then
        systemctl stop control-ofc-daemon.service
    fi
    if systemctl is-enabled --quiet control-ofc-daemon.service; then
        systemctl disable control-ofc-daemon.service
    fi
}

post_remove() {
    echo ":: Configuration preserved in /etc/control-ofc/"
    echo "   Remove manually if no longer needed."
    echo ":: State data preserved in /var/lib/control-ofc/"
    echo "   Remove manually if no longer needed."
}
